# backend/tradstry-backend.service
[Unit]
Description=Tradstry Backend API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tradstry-backend
ExecStart=/opt/tradstry-backend/tradstry-backend
Restart=always
RestartSec=10

# Environment variables (can be overridden in .env file)
Environment=RUST_LOG=info
Environment=PORT=9000

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tradstry-backend

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

